# 04/ex04.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{product}=="Apple Keyboard", RUN+="/sbin/insmod /root/little-pinguin/04/main.ko"
